# services/c3pUpgrade.service
[Unit]
Description=Cloud 3D Print Control Application Upgrade
After=multi-user.target

[Service]
WorkingDirectory=/home/pi
ExecStart=/home/pi/control/media/upgrade.sh
User=root
Type=simple
